# systemd/montessori-blog.timer
[Unit]
Description=Montessori Blog Automation Timer (diario)

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
